# cloudwatch-tools/roles/cw_memory/files/cloudwatch-memory.cron
* * * * * root /usr/local/bin/cw-memory-usage
