# cloudwatch-tools/roles/cw_disk/files/cloudwatch-disk.cron
* * * * * root /usr/local/bin/cw-disk-usage
